# web/httpd.conf
[program:httpd]
command=/usr/bin/pidproxy /var/run/httpd.pid /bin/bash -c "/usr/sbin/httpd -DFOREGROUND -k start"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisord/httpd-stderr.log
stdout_logfile=/var/log/supervisord/httpd-stdout.log